# makahiki/scripts/org.makahiki.evaluation.competition_summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>Makahiki: Check Energy Goal</string>
    <key>WorkingDirectory</key>
    <string>/Users/kukuicup/evaluation/makahiki</string>
    <key>ProgramArguments</key>
    <array>
    <string>/Users/kukuicup/Envs/kukuicup/bin/python</string>
    <string>competition_summary.py</string>
    <!-- Location to write out the file -->
    <string>/Users/kukuicup/evaluation/summary/summary.csv</string>
    </array>
    
    <!-- Only applies if the task is run as root. -->
    <key>UserName</key>
    <string>kukuicup</string>
    
    <!-- Run this task at :12 every hour. -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Minute</key>
      <integer>12</integer>
    </dict>
    
    <key>StandardErrorPath</key>
  	<string>/Users/kukuicup/evaluation/log/competition-summary.err</string>
  	<key>StandardOutPath</key>
  	<string>/Users/kukuicup/evaluation/log/competition-summary.out</string>
    	
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
